# flyscripts/files/etc/hotplug.d/button/01-button
#!/bin/sh

if [ $ACTION == "released" ]; then
  if [ $BUTTON == "BTN_2" ]; then
    logger "BTN_2 pressed"
  fi
fi
