# luci-theme-openwisp/root/etc/uci-defaults/31_luci-theme-openwisp
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Openwisp=/luci-static/openwisp
	set luci.main.mediaurlbase=/luci-static/openwisp
	commit luci
EOF
exit 0
