# luci-theme-octonix/root/etc/uci-defaults/31_luci-theme-octonix
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Octonix=/luci-static/octonix
	set luci.main.mediaurlbase=/luci-static/octonix
	commit luci
EOF
exit 0
